# contrib/crontab
*/2 * * * * /usr/local/bin/inventory send  >> /var/log/inventory.log 2>&1
